# tasks/0822639ce1f3deaf/seed/systemd/acme-logshipper.service
[Unit]
Description=ACME log shipping worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/acme-logshipper
Restart=on-failure
DynamicUser=yes
StateDirectory=acme-logshipper

[Install]
WantedBy=multi-user.target
